# scripts/agent-console.service.template
[Unit]
Description=Agent Console Server
After=network.target

[Service]
Type=simple
ExecStart={{COMMAND_PATH}} {{HOME}}/.agent-console/server/index.js
WorkingDirectory={{HOME}}/.agent-console/server
Restart=always
RestartSec=5

Environment=NODE_ENV=production
Environment=PORT={{PORT}}
Environment=APP_URL={{APP_URL}}
Environment=PATH={{PATH}}

StandardOutput=append:{{HOME}}/.local/share/agent-console/logs/agent-console.log
StandardError=append:{{HOME}}/.local/share/agent-console/logs/agent-console.error.log

[Install]
WantedBy=default.target
