# scripts/com.agent-console.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent-console</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{HOME}}/.agent-console/server/start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{HOME}}/.agent-console/server</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PORT</key>
        <string>{{PORT}}</string>
        <key>APP_URL</key>
        <string>{{APP_URL}}</string>
        <key>PATH</key>
        <string>{{PATH}}</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{HOME}}/Library/Logs/agent-console/agent-console.log</string>

    <key>StandardErrorPath</key>
    <string>{{HOME}}/Library/Logs/agent-console/agent-console.error.log</string>
</dict>
</plist>
